# ipopstat-0.12/DEBIAN/postinst
#!/usr/bin/env bash
service ipop-stat start
update-rc.d ipop-stat defaults
iptables -t nat -A PREROUTING -p tcp --dport 5000 -j REDIRECT --to-port 8080
