# ipopusagereport-0.1/DEBIAN/postinst
#!/usr/bin/env bash
service ipop-usagereport start
update-rc.d ipop-usagereport defaults
