# ipopstat-0.15/DEBIAN/postinst
#!/usr/bin/env bash
service ipop-stat start
update-rc.d ipop-stat defaults
